# config/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
username=worker
password=worker

[supervisord]
logfile=/home/worker/data/supervisor/logs/supervisord.log
pidfile=/tmp/supervisord.pid
logfile_maxbytes=50MB 
logfile_backups=10
loglevel=info
nodaemon=false
minfds=1024
minprocs=200
user=worker 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=/home/worker/php/sbin/php-fpm -c /home/worker/php/etc/php-fpm.ini
autostart=true
autorestart=true
priority=5
stdout_logfile=/home/worker/data/supervisor/logs/php-fpm-stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/home/worker/data/supervisor/logs/php-fpm-stderr.log
stderr_logfile_maxbytes=0
user=worker

[program:nginx]
command=/home/worker/nginx/sbin/nginx
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/home/worker/data/supervisor/logs/nginx-stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/home/worker/data/supervisor/logs/nginx-stderr.log
stderr_logfile_maxbytes=0
user=worker

[include]
files = /home/worker/supervisor/conf.d/*.conf